# deploy/systemd/contract-agent-api.service.template
[Unit]
Description=Contract Agent FastAPI Service
After=network-online.target ollama.service
Wants=network-online.target
Requires=ollama.service
# 避免 Restart=always 遇到啟動失敗時無限重啟打爆 CPU/journal。
# 60 秒內連續啟動失敗 5 次則 systemd 停止嘗試。
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__PROJECT_DIR__
Environment=HOME=__USER_HOME__
EnvironmentFile=__ENV_FILE__
ExecStartPre=/bin/bash -c 'fuser -k __API_PORT__/tcp 2>/dev/null || true'
ExecStart=/bin/bash -lc 'source __UV_ENV_FILE__ && cd __PROJECT_DIR__ && uv run uvicorn backend.main:app --host 0.0.0.0 --port __API_PORT__'
Restart=always
RestartSec=3
KillMode=control-group
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
