# deploy/systemd/contract-agent-web.service.template
[Unit]
Description=Contract Agent Vue Preview Service
After=network-online.target contract-agent-api.service
Wants=network-online.target
Requires=contract-agent-api.service

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__WEB_DIR__
Environment=HOME=__USER_HOME__
Environment=NODE_ENV=production
ExecStart=/bin/bash -lc 'cd __WEB_DIR__ && npm run preview -- --host 0.0.0.0 --port __WEB_PORT__'
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
