# deploy/systemd/contract-agent-web-admin.service.template
[Unit]
Description=Contract Agent Admin Vue Service
After=network-online.target contract-agent-api.service
Wants=network-online.target
Requires=contract-agent-api.service

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__WEB_DIR__
Environment=HOME=__USER_HOME__
Environment=NODE_ENV=production
ExecStartPre=/bin/bash -c 'fuser -k __WEB_ADMIN_PORT__/tcp 2>/dev/null || true'
ExecStart=/bin/bash -lc 'cd __WEB_DIR__ && npm run preview:admin -- --host 0.0.0.0 --port __WEB_ADMIN_PORT__ --strictPort'
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
